# core/scripts/crontab/drunken-robots.crontab
MAILTO="allmightyadmin@foo.bar"

# everyday at 22:05
5 22 * * * cd /home/pi/drunken-robots/robots; python manager.py | ../core/scr/core/scripts/mail/stdout_send.py

# TUDO: Get the right path
